# configuration_files/setup.conf
[program:container_setup]
autostart=true
autorestart=unexpected
redirect_stderr=true
user=root
directory=/setup_config
startsecs=0
command=/setup_config/setup-script.sh
