# configuration_files/supervisord-confs/apache-server-runner.conf
[program:apache-server-runner]
autostart=true
autorestart=true
redirect_stderr=true
user=root
startsecs=0
priority=999
command=httpd