# configuration_files/supervisord-confs/munge.conf
[program:munge]
command=munged -F
autostart=true
autorestart=true
redirect_stderr=true
user=munge
priority=10