# configuration_files/supervisord-confs/php-fpm-runner.conf
[program:php-fpm-server-runner]
autostart=false
autorestart=false
redirect_stderr=true
user=root
startsecs=0
priority=400
command=php-fpm -F